# systemd-timers/mail_recordatorio@.timer
[Unit]
Description=Enviar email de recordatorio todas las mañanas

[Timer]
OnCalendar=Mon,Fri *-*-* 7:30

[Install]
WantedBy=multi-user.target