# systemd-timers/mail_vencidos@.timer
[Unit]
Description=Enviar email de prestamos vencidos todas las mañanas

[Timer]
OnCalendar=Mon,Fri *-*-* 7:30

[Install]
WantedBy=multi-user.target