# systemd-timers/reindexar@.service
[Unit]
Description=Meran: reindexar en el contenedor %i

[Service]
Type=simple
ExecStart=/usr/bin/docker exec -d %i /cron/reindexar