# files/systemd-timers/generar_indice@.service
[Unit]
Description=Meran: generar indice en el contenedor %i

[Service]
Type=simple
ExecStart=/usr/bin/docker exec -d %i /cron/generar_indice