# files/systemd-timers/generar_indice@.timer
[Unit]
Description=Generar indice

[Timer]
OnCalendar=Mon *-*-* 6:00

[Install]
WantedBy=multi-user.target