# files/systemd-timers/portadas@.service
[Unit]
Description=Meran: descargar portadas en el contenedor %i

[Service]
Type=simple
ExecStart=/usr/bin/docker exec %i /cron/portadas